# ops/backup/day-log-backup.service
[Unit]
Description=Run Day Log MySQL backup
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/day-log
EnvironmentFile=/opt/day-log/.env
ExecStart=/usr/bin/docker compose --profile ops run --rm backup
