# ops/backup/daymark-backup.service
[Unit]
Description=Run Daymark MySQL backup
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/daymark
EnvironmentFile=/opt/daymark/.env
ExecStart=/usr/bin/docker compose --profile ops run --rm backup
